# scripts/com.soc.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.soc.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/soc-agent/soc-agent</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/soc-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/soc-agent-error.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/soc-agent</string>
</dict>
</plist>